# binaries/fpgad/fpgad.service.in
# Intel FPGA daemon
[Unit]
Description = FPGA monitor

[Service]
Type=simple
StandardOutput=journal+console
StandardError=journal+console
EnvironmentFile=-/etc/sysconfig/fpgad.conf
Restart=always
RestartSec=10s
SendSIGHUP=yes
KillSignal=SIGHUP
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/fpgad \
          -l $LOGFILE \
          -p $PIDFILE \
          -c $CFGFILE
RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target
